# deploy/systemd/ids-operator-console.service
[Unit]
Description=IDS operator console service
Wants=network-online.target
After=network-online.target

[Service]
Type=exec
User=ids-operator
Group=ids-operator
WorkingDirectory=/opt/ids_ml_new
EnvironmentFile=-/etc/ids-operator-console/ids-operator-console.env
ExecStartPre=/opt/ids_ml_new/.venv/bin/python -m ids.ops.operator_console_preflight --python-binary /opt/ids_ml_new/.venv/bin/python --app-module ids.console.server --manage-module ids.ops.operator_console_manage --database-path ${IDS_OPERATOR_CONSOLE_DATABASE_PATH} --alerts-input-path ${IDS_OPERATOR_CONSOLE_ALERTS_INPUT_PATH} --quarantine-input-path ${IDS_OPERATOR_CONSOLE_QUARANTINE_INPUT_PATH} --summary-input-path ${IDS_OPERATOR_CONSOLE_SUMMARY_INPUT_PATH} --templates-dir ${IDS_OPERATOR_CONSOLE_TEMPLATES_DIR} --static-dir ${IDS_OPERATOR_CONSOLE_STATIC_DIR} --environment ${IDS_OPERATOR_CONSOLE_ENVIRONMENT} --public-base-url ${IDS_OPERATOR_CONSOLE_PUBLIC_BASE_URL} --root-path ${IDS_OPERATOR_CONSOLE_ROOT_PATH} --forwarded-allow-ips ${IDS_OPERATOR_CONSOLE_FORWARDED_ALLOW_IPS} --secret-key-file ${IDS_OPERATOR_CONSOLE_SECRET_KEY_FILE}
ExecStart=/usr/bin/bash -lc 'exec /opt/ids_ml_new/.venv/bin/python -m ids.console.server --host "$IDS_OPERATOR_CONSOLE_HOST" --port "$IDS_OPERATOR_CONSOLE_PORT" --log-level "$IDS_OPERATOR_CONSOLE_LOG_LEVEL"'
Restart=on-failure
RestartSec=5s
StateDirectory=ids-operator-console
LogsDirectory=ids-operator-console
StandardOutput=journal
StandardError=journal
UMask=0027
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
